# incron.conf
# incron - event based background program processing daemon
#
# incron is a program that runs user-specified programs at
# filesystem events

description	"event based background program processing daemon"
author		"Andreas Altair Redmer <altair.ibn.la.ahad.sy@gmail.com>"

start on runlevel [2345]
stop on runlevel [!2345]

expect fork
respawn

exec incrond
